# scripts/orchid-serve.service
[Unit]
Description=Orchid persistent web server and agent orchestrator
After=network.target
Wants=network.target

[Service]
Type=simple
User=dave
WorkingDirectory=/home/dave/LocalAI/orchid
EnvironmentFile=/home/dave/LocalAI/orchid/.env
ExecStart=/home/dave/.local/bin/orchid serve \
    --watch-dir /home/dave/LocalAI \
    --watch-dir /home/dave/Documents/Development \
    --port 7842
Restart=on-failure
RestartSec=10
# Phase 1: send SIGTERM first (uvicorn graceful drain), then SIGKILL after 35s
# (5s buffer over runner.shutdown_timeout=30s)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=35
StandardOutput=journal
StandardError=journal

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/dave
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
